# iso/airootfs/etc/systemd/system/.service.d/liveuser.conf
[Unit]
Requires=arches-liveuser.service
After=arches-liveuser.service
